# calibre-web-rootless/root/etc/services.d/calibre-web/run
#!/usr/bin/with-contenv bash

cd /app/calibre-web || exit

exec \
    python3 /app/calibre-web/cps.py
